# .config/Code/User/History/-32213260/aDnq.service
[Unit]
Description=Battery Level Checker

[Service]
Restart=always
RestartSec=300
ExecStart=/home/yaros/Documents/battery_listener.sh