# .config/Code/User/History/-32213260/kY40.service
[Unit]
Description=Battery Level Checker

[Service]
ExecStart=/home/yaros/Documents/battery-check.sh