# .config/Code/User/History/-32213260/NPRB.service
[Unit]
Description=Battery Level Checker

[Service]
Restart=always
RestartSec=1
ExecStart=/home/yaros/Documents/battery_listener.sh