# .config/Code/User/History/-32213260/T5uJ.service
[Unit]
Description=Battery Level Checker

[Service]
Restart=always
RestartSec=5
ExecStart=/home/yaros/Documents/battery_listener.sh