# .config/Code/User/History/489ccf0/m8tN.timer
[Unit]
Description=Run Battery Level Checker every 5 minutes

[Timer]
OnBootSec=100
OnUnitActiveSec=100

[Install]
WantedBy=timers.target