# .config/Code/User/History/-59448d0/1gOD.service
[Unit]
Description=Battery Level Checker

[Service]
Restart=always
RestartSec=600
ExecStart=/home/yaros/Documents/battery_reset.sh