# .config/Code/User/History/3301ea18/kktZ.service
[Unit]
Description=Run Battery Level Checker every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target